# systemd/tailscale_exporter.socket
[Unit]
Description=Prometheus Tailscale Exporter

[Socket]
ListenStream=9824

[Install]
WantedBy=sockets.target
