# systemd/tailscale_exporter.timer
[Unit]
Description=Gather Tailscale metrics at regular interval

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
